# systemd/lacylights-gpio.service
[Unit]
Description=LacyLights GPIO Button Monitor
Documentation=https://github.com/bbernstein/lacylights-rpi
# Start after the main backend service is ready
After=lacylights.service
Requires=lacylights.service

[Service]
Type=simple
ExecStart=/opt/lacylights/scripts/gpio-button-monitor.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Limit restart attempts to prevent boot loops
StartLimitBurst=5
StartLimitIntervalSec=300

# Run as root for GPIO access
User=root
Group=root

# Environment
Environment="GPIO_PIN=17"
Environment="HOLD_TIME=5"
Environment="GRAPHQL_ENDPOINT=http://localhost:4000/graphql"

[Install]
WantedBy=multi-user.target
