# systemd/lacylights-wifi.service
[Unit]
Description=LacyLights WiFi Mode Management
Documentation=https://github.com/bbernstein/lacylights-rpi
# Wait for network and backend to be ready before checking WiFi mode
After=NetworkManager.service network-online.target lacylights.service
Wants=network-online.target
Requires=lacylights.service

[Service]
Type=oneshot
ExecStart=/opt/lacylights/scripts/wifi-mode-check.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

# Timeout to prevent boot delays if script hangs
TimeoutStartSec=60s

# Environment
Environment="WIFI_CHECK_TIMEOUT=30"
Environment="WIFI_INTERFACE=wlan0"
Environment="GRAPHQL_ENDPOINT=http://localhost:4000/graphql"

[Install]
WantedBy=multi-user.target
