# systemd/lacylights-frontend.service
[Unit]
Description=LacyLights Frontend (Next.js Server)
# Wait for network to be fully online before starting
# This ensures the network interface has an IP address
After=network-online.target
Wants=network-online.target lacylights.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/lacylights/frontend-src
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="NEXT_PUBLIC_GRAPHQL_ENDPOINT=http://lacylights.local/graphql"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lacylights-frontend

[Install]
WantedBy=multi-user.target
