# systemd/lacylights.service
[Unit]
Description=LacyLights Stage Lighting Control (Go Backend)
Documentation=https://github.com/bbernstein/lacylights
# Wait for network to be fully online before starting
# This ensures the network interface has an IP address
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lacylights
Group=lacylights
WorkingDirectory=/opt/lacylights/backend

# Environment variables (optional - continues if file doesn't exist)
EnvironmentFile=-/opt/lacylights/backend/.env

# Start the Go backend server
ExecStart=/opt/lacylights/backend/lacylights-server

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=3
StartLimitInterval=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lacylights

# Security settings
# NoNewPrivileges disabled to allow self-updates via sudo systemd-run
# The sudoers file restricts which commands can be run with elevated privileges
PrivateTmp=true
ProtectSystem=strict
# ProtectHome=read-only allows access to /home but not modification (except via ReadWritePaths)
ProtectHome=read-only
# Allow write access to:
#   - /opt/lacylights: deployment directories
#   - /home/lacylights: npm cache for updates
#   - /run/sudo: required for sudo timestamp files (sudo needs to cache auth tokens)
ReadWritePaths=/opt/lacylights /home/lacylights /run/sudo

# Resource limits
LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target
